# ansible/roles/leg_service/service/leg.service
[Unit]
Description=LEG Script Service

[Service]
ExecStart=/usr/local/bin/leg.sh
